# tailscale/rootfs/etc/s6-overlay/s6-rc.d/taildrop/run
#!/command/with-contenv bashio
# shellcheck shell=bash
export LOG_FD
# ==============================================================================
# Home Assistant Community App: Tailscale
# Fetches files send via Taildrop
# ==============================================================================

# Add timestamp to plain stdout and stderr messages
exec &> >(ts "${__BASHIO_LOG_TIMESTAMP} service $(basename ${PWD}):")

# Ensure the directory exists
mkdir -p /share/taildrop

# Fetch files
exec /opt/tailscale file get --verbose --loop "/share/taildrop"
