# tailscale/rootfs/etc/s6-overlay/s6-rc.d/init-protect-subnets/run
#!/command/with-contenv bashio
# shellcheck shell=bash
export LOG_FD
# ==============================================================================
# Home Assistant Community App: Tailscale
# Prevent local subnets to be routed toward the tailnet
# ==============================================================================

# Add timestamp to plain stdout and stderr messages
exec &> >(ts "${__BASHIO_LOG_TIMESTAMP} service $(basename ${PWD}):")

protect-subnet-routes
