# tailscale/rootfs/etc/s6-overlay/s6-rc.d/init-magicdns-proxies/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
export LOG_FD
# ==============================================================================
# Home Assistant Community App: Tailscale
# Remove forwarding
# ==============================================================================

# Add timestamp to plain stdout and stderr messages
exec &> >(ts "${__BASHIO_LOG_TIMESTAMP} service $(basename ${PWD}):")

magicdns-ingress-proxy-forwarding remove drop
