# tailscale/rootfs/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
export LOG_FD
# ==============================================================================
# Home Assistant Community App: Tailscale
# Configures NGINX for use with the Tailscale web
# ==============================================================================

# Add timestamp to plain stdout and stderr messages
exec &> >(ts "${__BASHIO_LOG_TIMESTAMP} service $(basename ${PWD}):")

# Generate Ingress configuration
bashio::var.json \
    interface "$(bashio::addon.ip_address)" \
    port "^$(bashio::addon.ingress_port)" \
    | tempio \
        -template /etc/nginx/templates/ingress.gtpl \
        -out /etc/nginx/servers/ingress.conf
